# puppet/modules/uwsgi/templates/uwsgi.conf.erb
# Emperor uWSGI script

description "uWSGI Emperor"
start on vagrant-mounted MOUNTPOINT=<%= @mount_point %>
stop on runlevel [06]

exec /usr/local/bin/uwsgi <% @params.each_pair do | param, val | %> --<%= param %> <%= val %> <% end %>
